# linux-client/accutime-sync.service
[Unit]
Description=Set system clock from AccuTime phone GPS over WiFi hotspot (NTP)
After=network-online.target
Wants=network-online.target

[Service]
# Keeps the clock disciplined while connected to the phone's hotspot.
# 'auto' detects the phone as the default gateway — no IP to edit for field use.
Type=simple
ExecStart=/usr/bin/python /usr/local/bin/accutime-sync.py auto --loop --interval 300 --hwclock
Restart=on-failure
RestartSec=30
# systemd 229 (Ubuntu 16.04 / Linux Lite 3.4) accepts these names.
StartLimitInterval=0

[Install]
WantedBy=multi-user.target
